# gunicorn.conf
[program:gunicorn]
command=gunicorn --pythonpath . -k gevent app:app -b 0.0.0.0:5000
directory=/deploy/flask_rest_template
